# iql/crontab
*/1  *  *  *  * shardbuilder /opt/imhotepTsvConverter/tsvConverter.sh
