# iql/crontab
*/1  *  *  *  * shardbuilder /opt/imhotepShardBuilder/shardBuilder.sh
